# conf/auto_sync.service
[Unit]
Description=auto_sync (scheduler + watcher + web)
After=local-fs.target network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/auto_sync
ExecStart=/opt/auto_sync/bin/auto_sync --config /opt/auto_sync/conf/auto_sync.toml
Restart=always
RestartSec=5
User=root
Group=root
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_SYS_ADMIN CAP_DAC_READ_SEARCH
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
